# deploy/supervisor/discord-bot.conf
[program:discord-bot]
command=/usr/bin/python3 /opt/discord-bot-ai/main.py
directory=/opt/discord-bot-ai
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30

# Run as dedicated user
user=discordbot
group=discordbot

# Environment
environment=HOME="/opt/discord-bot-ai",USER="discordbot"
environmentfile=/opt/discord-bot-ai/.env.docker

# Logging
stdout_logfile=/var/log/discord-bot/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile_makepartial=true
stderr_logfile=/var/log/discord-bot/stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

# Process management
stopasgroup=true
killasgroup=true
priority=999

# Security
no-new-privileges=true

# Resource limits (supervisor 4.0+)
# Acceptable limits: RLIMIT_AS (memory), RLIMIT_CPU (cpu time)
# Note: These are soft limits
# limit_as=512M
